# root/etc/s6-overlay/s6-rc.d/svc-mod-qbittorrent-port/run
#!/usr/bin/with-contenv bash

QBIT_HOST="${QBIT_HOST:-localhost}"
QBIT_PORT="${QBIT_PORT:-8080}"
QBIT_USER="${QBIT_USER:-admin}"
QBIT_PASS="${QBIT_PASS:-adminadmin}"
QBIT_URL="http://${QBIT_HOST}:${QBIT_PORT}"

echo "Waiting for qBittorrent to be ready..."
until curl -s -o /dev/null "${QBIT_URL}/api/v2/app/version"; do
    sleep 5
done
echo "qBittorrent is ready"

while true; do
    gluetun_response=$(curl -s localhost:8000/v1/openvpn/portforwarded)
    if [ $? -ne 0 ]; then
        echo "Failed to reach Gluetun API"
        sleep 10
        continue
    fi

    gluetun_port=$(echo "$gluetun_response" | jq -r '.port')

    if [ -z "$gluetun_port" ] || [ "$gluetun_port" == "null" ] || [ "$gluetun_port" == "0" ]; then
        echo "No valid port from Gluetun (got: $gluetun_port)"
        sleep 10
        continue
    fi

    # Login to qBittorrent
    cookie=$(curl -s -i --header "Referer: ${QBIT_URL}" --data "username=${QBIT_USER}&password=${QBIT_PASS}" "${QBIT_URL}/api/v2/auth/login" 2>&1)
    if [ $? -ne 0 ]; then
        echo "Failed to login to qBittorrent"
        sleep 10
        continue
    fi

    sid=$(echo "$cookie" | grep -i 'set-cookie' | sed 's/.*SID=\([^;]*\).*/\1/')

    if [ -z "$sid" ]; then
        echo "Failed to get session cookie from qBittorrent, raw response:"
        echo "$cookie"
        sleep 10
        continue
    fi

    # Get current listen port
    prefs=$(curl -s --cookie "SID=${sid}" "${QBIT_URL}/api/v2/app/preferences" 2>&1)
    if [ $? -ne 0 ]; then
        echo "Failed to get qBittorrent preferences: $prefs"
        sleep 10
        continue
    fi

    current_port=$(echo "$prefs" | jq -r '.listen_port')

    if [ -z "$current_port" ]; then
        echo "Failed to parse qBittorrent listen port, raw output:"
        echo "$prefs"
        sleep 10
        continue
    fi

    if [ "$current_port" == "$gluetun_port" ]; then
        echo "qBittorrent port is already set to $gluetun_port"
        sleep 60
        continue
    fi

    # Set new listen port
    update_output=$(curl -s --cookie "SID=${sid}" --data "json={\"listen_port\": ${gluetun_port}}" "${QBIT_URL}/api/v2/app/setPreferences" 2>&1)
    if [ $? -ne 0 ]; then
        echo "Failed to update qBittorrent port: $update_output"
        sleep 10
        continue
    fi

    # Reannounce all torrents
    curl -s --cookie "SID=${sid}" --data "hashes=all" "${QBIT_URL}/api/v2/torrents/reannounce" >/dev/null 2>&1

    echo "qBittorrent port updated to $gluetun_port"

    sleep 60
done
